# tools/docker/agent-instance/inittab
::sysinit:bash /etc/init.d/agent-instance-startup
::respawn:/usr/bin/monit -Ic /etc/monitrc
#::askfirst:-/bin/sh
