# resources/content/config-content/healthcheck/content/etc/init.d/haproxy-monitor
#!/sbin/openrc-run

CATTLE_HOME=${CATTLE_HOME:-/var/lib/cattle}
PIDFILE=/var/run/${SVCNAME}.pid
DAEMON=${CATTLE_HOME}/bin/host-api
DAEMON_ARGS="-log /var/log/haproxy-monitor.log -haproxy-monitor -pid-file ${PIDFILE}"

start() {
	ebegin "Starting ${SVCNAME}"
	mkdir -p $CATTLE_HOME
	start-stop-daemon --start -b -d $CATTLE_HOME --pidfile $PIDFILE --exec $DAEMON -- $DAEMON_ARGS
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --exec $DAEMON --pidfile $PIDFILE
	rm -f $PIDFILE
	eend $?
}
